# pi/systemd/sllidar.service
[Unit]
Description=RPLIDAR A2M12 laser scan node
After=multi-user.target

[Service]
User=pi
# 注意:務必用 /dev/rplidar 短別名(sllidar SDK 有埠名長度 bug)
ExecStart=/bin/bash -c 'source /opt/ros/jazzy/setup.bash && source /home/pi/ros2_ws/install/setup.bash && exec ros2 run sllidar_ros2 sllidar_node --ros-args -p serial_port:=/dev/rplidar -p serial_baudrate:=256000 -p frame_id:=laser'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
